# powershell/profile.d/env.ps1
# Shared PowerShell environment defaults.
if (-not $env:OP_ACCOUNT) {
  $env:OP_ACCOUNT = "CNRNCJQSBBBYZESUWAMXLHQFBI"
}
