# powershell/profile.d/op.ps1
# Default 1Password account for the `op` CLI.
if (-not $env:OP_ACCOUNT) {
  $env:OP_ACCOUNT = "CNRNCJQSBBBYZESUWAMXLHQFBI"
}
