# systemd/user/openclaw-gateway.service.d/override.conf
[Service]
EnvironmentFile=%h/.openclaw/gateway.env
