# systemd/user/openclaw-gateway.service.d/override.conf
[Service]
# Optional: Load gateway.env if exists (missing file won't fail service startup)
EnvironmentFile=-%h/.openclaw/gateway.env
