# systemd/user/openclaw-cleanup.timer
[Unit]
Description=Run cleanup at ~05:00 daily
After=default.target

[Timer]
OnCalendar=*-*-* 05:00:00
Persistent=true
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target
