# systemd/user/openclaw-cleanup.service
[Unit]
Description=Periodic cleanup (mise/pnpm/npm) to reduce disk usage

[Service]
Type=oneshot
# Run the script deployed by Home Manager to ~/.config/scripts/openclaw-cleanup
ExecStart=%h/.config/scripts/openclaw-cleanup
