# systemd/user/openclaw-cleanup.service
[Unit]
Description=Periodic cleanup (mise/pnpm/npm) to reduce disk usage

[Service]
Type=oneshot
# Installed via dotfiles into ~/bin (see zsh/.zprofile PATH).
ExecStart=%h/bin/openclaw-cleanup
